# server/start-script/server-startup.service
[Unit]
Description: Service to Launch Fastapi server

[Service]
ExecStart=/usr/bin/server-startup.sh

[Install]
WantedBy=multi-user.target
